# .config/systemd/user/pulseaudio.service
[Unit]
Description=pulseaudio daemon

[Service]
ExecStart=pulseaudio

[Install]
WantedBy=default.target
