# .config/systemd/user/udiskie.service
[Unit]
Description=udiskie - automount USB disks

[Service]
ExecStart=udiskie -F

[Install]
WantedBy=default.target
